# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:recto]
command=/bin/bash -c "cd /opt/sinopia/recto/ && npm start &"

[program:verso]
command=/bin/bash -c "cd /opt/sinopia/verso/ && npm start &"
